# ex.conf/pingHosts.conf
[program:pingMonitor]
command=/home/gx/Documents/projects/lcd-monitor/pingHosts.py
process_name=pingMonitor
autostart=true
autorestart=true

